# etc/supervisord.d/redis.ini
[program:redis]
command=redis-server
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
autorestart = true
priority = 10
