# etc/supervisord.d/sentry-web.ini
[program:sentry-web]
command=su -c 'sentry start' user
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
autorestart = true
priority = 10
