# etc/supervisord.d/sentry-worker.ini
[program:sentry-worker]
command=/opt/qnib/sentry/bin/start_celery_workers.sh
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
autorestart = true
priority = 10
