# etc/supervisord.d/sentry-beat.ini
[program:sentry-beat]
command=/opt/qnib/sentry/bin/start_celerybeat.sh
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
autorestart = true
priority = 10
